# scripts/com.beamr.norn-pr-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.beamr.norn-pr-review</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/joshua/Developer/beamr/scripts/norn-pr-review.sh</string>
        <string>--poll</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/joshua/Developer/beamr</string>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>StandardOutPath</key>
    <string>/Users/joshua/Library/Logs/norn-pr-review.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/joshua/Library/Logs/norn-pr-review.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/joshua/.cargo/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/joshua</string>
    </dict>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
